# .kairos/k-chat-watchdog.service
[Unit]
Description=K-Chat Watchdog — self-healing crash recovery daemon
After=network-online.target k-chat.service
Wants=network-online.target k-chat.service

[Service]
Type=simple
WorkingDirectory=%h/dev/K-Chat

# ── Environment ────────────────────────────────────────────────────
Environment=WATCHDOG_INTERVAL=15
Environment=WATCHDOG_REQUIRED_FAILURES=6
Environment=WATCHDOG_URL=http://127.0.0.1:8000/health
Environment=WATCHDOG_STARTUP_TIMEOUT=180
Environment=KAIROS_WEB_SERVICE=k-chat.service
Environment=PYTHONUNBUFFERED=1

# ── Execute ────────────────────────────────────────────────────────
# The watchdog only monitors health and asks k-chat.service to restart.
# systemd remains the single owner of the uvicorn process.
ExecStart=%h/dev/K-Chat/venv/bin/python .kairos/watchdog.py

# ── Lifecycle ──────────────────────────────────────────────────────
Restart=on-failure
RestartSec=5

# ── Security ───────────────────────────────────────────────────────
User=%u
NoNewPrivileges=yes
PrivateTmp=yes
# The watchdog writes crash context and cache invalidation markers.
ProtectHome=read-only
ProtectSystem=full
ReadWritePaths=%h/dev/K-Chat

[Install]
WantedBy=default.target
